# deploy/logoust-assistant.service
[Unit]
Description=Logoust Assistant Telegram bot (long-polling)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=logoust
Group=logoust
WorkingDirectory=/opt/logoust-assistant
EnvironmentFile=/opt/logoust-assistant/.env
# uv lives in the logoust user's home; full path so PATH need not be set.
ExecStart=/home/logoust/.local/bin/uv run --no-dev python -m src
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
